# systemd/template.service
[Unit]
Description=Template description

[Service]
Type=simple
User=root
WorkingDirectory=
ExecStart=
Restart=
RestartSec=

[Install]
WantedBy=